# deploy/systemd/world-observer-wiesmoor-weather.timer
[Unit]
Description=Run World Observer Wiesmoor weather refresh hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=world-observer-wiesmoor-weather.service

[Install]
WantedBy=timers.target
