# deploy/systemd/world-observer-iss-tle.service
[Unit]
Description=World Observer ISS TLE cache refresh
Documentation=https://github.com/dennishilk/world-observer
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=nebu
WorkingDirectory=/opt/world-observer/world-observer
ExecStart=/opt/world-observer/world-observer/scripts/update_iss_tle.py --path /opt/world-observer/world-observer/data/reference/iss.tle
