# cron/daily.cron
# world-observer daily automation (user: nebu)

# Heartbeat: hourly
0 * * * * /bin/bash -lc 'cd /home/nebu/world-observer && source .venv/bin/activate && python scripts/heartbeat_push.py >> /home/nebu/world-observer/logs/heartbeat.log 2>&1'

# Daily observer run: at 02:00 UTC
0 2 * * * /bin/bash -lc 'cd /home/nebu/world-observer && source .venv/bin/activate && python scripts/run_daily_cron.py >> /home/nebu/world-observer/logs/daily.log 2>&1'
